# conf/systemd_check.timer
[Unit]
Description=Check __APP__ backup regularly

[Timer]
OnCalendar=__CHECK_ON_CALENDAR__

[Install]
WantedBy=timers.target
